# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -nolisten tcp
priority=10
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:fluxbox]
command=/bin/bash -c "until xdpyinfo -display :99 >/dev/null 2>&1; do sleep 0.5; done && exec fluxbox"
environment=DISPLAY=":99"
priority=20
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:x11vnc]
command=/bin/bash -c "until xdpyinfo -display :99 >/dev/null 2>&1; do sleep 0.5; done && exec x11vnc -display :99 -forever -shared -nopw -rfbport 5900 -quiet"
priority=30
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:novnc]
command=/usr/share/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080
priority=40
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:uvicorn]
command=uv run uvicorn agentic_fetch.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
priority=50
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
